# packages/a/atftp/atftpd.socket
[Unit]
Description=Advanced tftp Server Activation Socket

[Socket]
ListenDatagram=0.0.0.0:69

[Install]
WantedBy=sockets.target
