# packages/c/coredns-for-k8s1.31/coredns.service
[Unit]
Description=Fast and flexible DNS server
After=network-online.target
Wants=network-online.target
Documentation=man:coredns(1)

[Service]
Type=simple
ExecStart=/usr/sbin/coredns -conf /etc/coredns/Corefile -pidfile /var/run/coredns.pid

[Install]
WantedBy=multi-user.target
