# packages/v/virtualbox/vboxservice.service
[Unit]
Description=VirtualBox guest services
ConditionVirtualization=|oracle

[Service]
Type=simple
ExecStartPre=-/usr/sbin/modprobe vboxguest
ExecStart=/usr/sbin/VBoxService -f
Restart=on-failure
