# packages/v/virtualbox/vboxclient.service
[Unit]
Description=VirtualBox guest VMSVGA resize client
ConditionVirtualization=|oracle
ExecCondition=sh -c '[[ "$XDG_SESSION_TYPE" == "wayland" ]] || exit -1'

[Service]
Type=simple
ExecStart=/usr/bin/VBoxDRMClient
Restart=on-failure

[Install]
WantedBy=multi-user.target

