# packages/v/virtualbox/vboxclient.service
[Unit]
Description=VirtualBox guest VMSVGA resize client
ConditionVirtualization=|oracle

[Service]
Type=simple
ExecStart=/usr/bin/VBoxDRMClient
Restart=on-failure

[Install]
WantedBy=multi-user.target
# also Wanted by 90-vboxguest.rules
