# packages/s/syslog-ng/syslog-ng.service
[Unit]
Description=System Logging Service
Requires=var-run.mount
# Requires=syslog.socket
After=var-run.mount network.target
Conflicts=rsyslog.service syslogd.service

[Service]
Environment=SYSLOG_NG_PARAMS=
ExecStartPre=@sbindir@/syslog-ng-service-prepare
EnvironmentFile=-/etc/sysconfig/syslog
ExecStart=@sbindir@/syslog-ng -F $SYSLOG_NG_PARAMS
ExecReload=/bin/kill -HUP $MAINPID
StandardOutput=null
Restart=on-abort

[Install]
WantedBy=multi-user.target
Alias=syslog.service
