# packages/s/suse-build-key/suse-build-key-import.timer
[Unit]
Description=Timer to import new SUSE RPM signing keys

[Timer]
OnBootSec=1h
OnCalendar=hourly

[Install]
WantedBy=timers.target
