# packages/d/dns24-update/dns24@.timer
[Unit]
Description=Run %p@%i every 22 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=22min

[Install]
WantedBy=default.target
