# packages/t/tftp/tftp.service
[Unit]
Description=Tftp Server
Wants=network.target
After=network.target
Requires=tftp.socket

[Service]
EnvironmentFile=/etc/sysconfig/tftp
ExecStart=/usr/sbin/in.tftpd -u $TFTP_USER -s $TFTP_OPTIONS $TFTP_DIRECTORY
StandardInput=socket
PrivateDevices=yes
