# packages/o/open-vm-tools/vmblock-fuse.service
[Unit]
Description=open-vm-tools: vmblock-fuse service for virtual machines hosted on VMware
Documentation=https://github.com/vmware/open-vm-tools/blob/master/open-vm-tools/vmblock-fuse/design.txt
ConditionVirtualization=vmware
Before=vmtoolsd.service

[Service]
Type=forking
RuntimeDirectory=vmblock-fuse
RuntimeDirectoryMode=755
ExecStart=/usr/bin/vmware-vmblock-fuse -o subtype=vmware-vmblock,default_permissions,allow_other /var/run/vmblock-fuse

[Install]
WantedBy=multi-user.target
