# packages/o/open-vm-tools/vmtoolsd.service
[Unit]
Description=open-vm-tools: vmtoolsd service for virtual machines hosted on VMware
Documentation=http://github.com/vmware/open-vm-tools
ConditionVirtualization=vmware
Requires=vgauthd.service
After=vgauthd.service
DefaultDependencies=no
Before=cloud-init-local.service

[Service]
Type=forking
PIDFile=/var/run/vmtoolsd.pid
ExecStart=/usr/bin/vmtoolsd --background /var/run/vmtoolsd.pid
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
Also=vgauthd.service
