# packages/g/git/git-daemon.service
[Unit]
Description=Start Git Daemon

[Service]
EnvironmentFile=-/etc/sysconfig/git-daemon
ExecStart=/bin/sh -c 'exec git daemon --reuseaddr --base-path=${GIT_DAEMON_BASE_PATH:-/srv/git/} --user=${GIT_DAEMON_USER:-git-daemon} --group=${GIT_DAEMON_GROUP:-nogroup} $GIT_DAEMON_ARGS'

Restart=always
RestartSec=500ms

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=git-daemon

[Install]
WantedBy=multi-user.target
