# packages/f/fwnn/fkwnn.service
[Unit]
ConditionPathExists=/etc/FreeWnn/ko_KR/kserverrc
Description=Free Wnn Korean Server, for input of Korean

[Service]
Type=forking
User=wnn
ExecStartPre=/bin/rm -f /tmp/kd_sockV4
ExecStart=/usr/bin/kserver
ExecStopPost=/bin/rm -f /tmp/kd_sockV4

[Install]
WantedBy=multi-user.target
