# packages/p/python-mailman/mailman-digests.service
[Unit]
Description=Sends Mailman digests

[Service]
Type=oneshot
ExecStart=/usr/bin/mailman digests --send
User=mailman
Group=mailman
