# packages/p/python-mailman/mailman-notify.timer
[Unit]
Description=Run mailman-notify.service every day

[Timer]
OnCalendar=*-*-* 08:00:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
