# ulogd/ulogd.service
[Unit]
Description=Userspace logging daemon for NFLOG

[Service]
ExecStart=/usr/sbin/ulogd -u ulogd

[Install]
WantedBy=multi-user.target
