# dns24-update/dns24@.service
[Unit]
Description=Update dynamic DNS address of '%i'

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -s --config /etc/dns24/%i.curl
