# packages/w/wireless-tools/99-iwlwifi-led.rules
SUBSYSTEM=="leds", ACTION=="add", KERNEL=="iwl-phy*:assoc", RUN+="@UDEVDIR@/iwlwifi-led.sh"

