# packages/t/tmpwatch/tmpwatch.timer
[Unit]
Description=Automatic clean-up execution of tmpwatch

[Timer]
OnCalendar=Daily
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
