# packages/n/nextcloud/nextcloud-cron.service
[Unit]
Description=Nextcloud CronJob

[Service]
User=wwwrun
ExecStart=/usr/bin/php -f /srv/www/htdocs/nextcloud/cron.php

[Install]
WantedBy=basic.target
