# packages/n/nextcloud/nextcloud-cron.service
[Unit]
Description=Nextcloud CronJob

[Service]
User=@APACHE_USER@
ExecCondition=/usr/bin/php -f @APACHE_MYSERVERROOT@/nextcloud/occ status -e
ExecStart=/usr/bin/php -f @APACHE_MYSERVERROOT@/nextcloud/cron.php
KillMode=process
