# packages/n/nextcloud/nextcloud-cron.service
[Unit]
Description=Nextcloud CronJob

[Service]
User=@APACHE_USER@
ExecStart=/usr/bin/php -f @APACHE_MYSERVERROOT@/nextcloud/cron.php

[Install]
WantedBy=basic.target
