# packages/s/sarg/sarg-reports-monthly.timer
[Unit]
Description=Monthly execution of sarg-reports

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target