# packages/s/sarg/sarg-reports-daily.timer
[Unit]
Description=Daily execution of sarg-reports

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target