# packages/b/bird/bird6.service
[Unit]
Description=The BIRD Internet Routing Daemon (IPv6)
Documentation=http://bird.network.cz/doc/bird.html
Wants=network.target
After=network.target
ConditionFileIsExecutable=/usr/sbin/bird6

[Service]
Type=simple
ExecStartPre=/usr/sbin/bird6 -p -u bird -g bird
ExecStart=/usr/sbin/bird6 -u bird -g bird -f -P /run/bird/bird6.pid $BIRD_OPTIONS
ExecReload=/usr/bin/kill -HUP $MAINPID
PIDFile=/run/bird/bird6.pid
Restart=on-failure

[Install]
WantedBy=multi-user.target
