# packages/libv/libvma/vma.service
[Unit]
Description=VMA Monitoring Daemon
StopWhenUnneeded=no
DefaultDependencies=no
# Ensure required kernel modules are loaded before starting
Wants=rdma-load-modules@rdma.service
After=rdma-load-modules@rdma.service
# After basic RDMA hw setup.
After=rdma-hw.target

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/sbin/vmad

[Install]
WantedBy=multi-user.target

