# packages/z/zabbix/zabbix-server.service
[Unit]
Description=Zabbix Server
After=syslog.target network.target

[Service]
Type=simple
ExecReload=/usr/sbin/zabbix-server -R config_cache_reload
ExecStart=/usr/sbin/zabbix-server -f
User=zabbixs
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
