# packages/z/zabbix/zabbix-agentd.service
[Unit]
Description=Zabbix Monitor Agent
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/usr/sbin/zabbix-agentd -f
User=zabbix
Group=zabbix
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
