# packages/l/leafnode/leafnode.socket
[Unit]
Description=Leafnode NNTP server

[Socket]
ListenStream=119
Accept=true

[Install]
WantedBy=sockets.target
