# packages/v/vsftpd/vsftpd.service
[Unit]
Description=Vsftpd ftp daemon
After=network.target

[Service]
ExecStart=/usr/sbin/vsftpd /etc/vsftpd.conf

[Install]
WantedBy=multi-user.target
