# packages/p/python-mautrix-telegram/mautrix-telegram.service
[Unit]
Description=A Matrix-Telegram hybrid puppeting/relaybot bridge.

[Service]
Type=simple
SyslogIdentifier=mautrix-telegram
Restart=on-failure
RestartSec=3

User=root
WorkingDirectory=/usr/lib/python3.*/site-packages//mautrix-telegram
ExecStart=@PYTHON_FLAVOR@ -m mautrix_telegram -c /etc/mautrix-telegram/config.yaml

[Install]
WantedBy=multi-user.target
