# packages/l/live-langset-data/langset.service
[Unit]
Description=Adjust the system language based on the kernel cmdline
Requires=dbus.service sysinit.target
Before=multi-user.target shutdown.target
After=dbus.service sysinit.target
ConditionPathExists=!/var/lib/zypp/RequestedLocales
 
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh /usr/sbin/langset.sh

[Install]
WantedBy=multi-user.target
