# erlang/epmd.service
[Unit]
Description=Erlang Port Mapper Daemon
After=network.target
Requires=epmd.socket

[Service]
ExecStart=/usr/bin/epmd -systemd
#ExecStop=/usr/bin/epmd -kill
Type=simple
StandardOutput=journal
StandardError=journal
Restart=always
LimitNPROC=1
LimitFSIZE=0
User=epmd
Group=epmd

[Install]
Also=epmd.socket
WantedBy=multi-user.target

