# container-feeder/container-feeder.service
[Unit]
Description=Load all docker images that are packaged in RPM
After=docker.service
Before=kubelet.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
EnvironmentFile=-/etc/sysconfig/container-feeder
ExecStart=/usr/bin/container-feeder $OPTS

[Install]
WantedBy=multi-user.target
