# packages/b/backupninja/backupninja.service
[Unit]
Description=Run backupninja
Requires=default.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/backupninja
