# packages/d/distribution/registry.service
[Unit]
Description=Docker Registry Server
Documentation=https://github.com/docker/distribution
After=network.target

[Service]
ExecStart=/usr/bin/registry serve /etc/registry/config.yml
User=registry
Group=registry

[Install]
WantedBy=multi-user.target

