# packages/f/flatpak/update-system-flatpaks.service
[Unit]
Description=Update system Flatpaks
Documentation=man:flatpak-update(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak --system update -y --noninteractive

[Install]
WantedBy=default.target
