# packages/p/python-HyperKitty/hyperkitty-runjob.service
[Unit]
Description=HyperKitty runjob @HYPERKITTY_RUNJOB@

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /srv/www/webapps/mailman/hyperkitty/manage.py runjob @HYPERKITTY_RUNJOB@ --pythonpath /srv/www/webapps/mailman/hyperkitty/ --settings settings
User=hyperkitty
Group=hyperkitty
