# packages/s/snowflake/snowflake.service
[Unit]
Description=Snowflake proxy
After=network.target

[Service]
Restart=on-failure
ExecStart=/usr/sbin/snowflake

[Install]
WantedBy=multi-user.target
