# packages/c/cfengine/cf-execd.service
[Unit]
Description=CFEngine Execution Daemon
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/bin/cf-execd

[Install]
WantedBy=multi-user.target
