# packages/v/valkey/valkey-sentinel.target
[Unit]
Description=Valkey Sentinel target allowing to start/stop all valkey-sentinel@.service instances at once

