# packages/v/valkey/valkey.target
[Unit]
Description=Valkey target allowing to start/stop all valkey@.service instances at once

