# packages/t/tayga/tayga.service
[Unit]
Description=Simple, no-fuss NAT64 for Linux
After=syslog.target network.target firewall.target

[Service]
Type=forking
ExecStartPre=/usr/sbin/tayga_setup_tun
ExecStart=/usr/sbin/tayga

[Install]
WantedBy=multi-user.target
