# packages/k/killerd/killerd.service
[Unit]
Description=Daemon for automatic killing of login shells
After=local-fs.target

[Service]
ExecStart=/usr/sbin/killerd

[Install]
WantedBy=multi-user.target
