# packages/v/varnish/varnish.service
[Unit]
Description=Varnish HTTP accelerator/reverse proxy
After=network.target

[Service]
EnvironmentFile=/etc/sysconfig/varnish
PIDFile=/var/run/varnishd.pid
ExecStart=/usr/sbin/varnishd -P /var/run/varnishd.pid -F $VARNISHD_PARAMS

[Install]
WantedBy=multi-user.target
