# packages/y/yum/yum-updatesd.service
[Unit]
Description=YUM Package	Update Service

[Service]
Type=dbus
BusName=edu.duke.linux.yum
ExecStart=/usr/sbin/yum-updatesd -f

[Install]
WantedBy=multi-user.target
