# packages/n/nagios/nagios-archive.timer
[Unit]
Description=Auto-Archiving of Nagios Logfiles

[Timer]
OnCalendar=Sat *-*-* 03:00:00
AccuracySec=12h
Persistent=true

[Install]
WantedBy=timers.target
