# packages/n/nagios/nagios-archive.service
[Unit]
Description=Auto-Archiving of Nagios Logfiles

[Service]
Type=oneshot
ExecStart=/bin/bash -ce "/usr/sbin/nagios-archive"
