# packages/n/nagios/nagios.service
[Unit]
Description=Network Monitor Nagios
After=syslog.target network.target

[Service]
EnvironmentFile=/etc/sysconfig/nagios
ExecStartPre=/usr/lib/nagios/nagios-exec-start-pre
ExecStartPre=/usr/sbin/nagios -v /etc/nagios/nagios.cfg
ExecStart=/usr/sbin/nagios /etc/nagios/nagios.cfg
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
Alias=monitoring_daemon.service
