# packages/g/gkrellm/gkrellmd.service
[Unit]
Description=GNU Krell Monitors server
After=network.target

[Service]
ExecStart=/usr/bin/gkrellmd

[Install]
WantedBy=multi-user.target

