# packages/q/quota/quotad.service
[Unit]
Description=Quota RPC monitor
After=network.target

[Service]
EnvironmentFile=-/run/sysconfig/quotad
ExecStart=/usr/sbin/rpc.rquotad -F ${RQUOTAD_ARGS}
ExecStartPre=-/usr/lib/systemd/scripts/quotad_env.sh

[Install]
WantedBy=multi-user.target
