# packages/c/clair/clair.service
[Unit]
Description=Clair Vulnerability Scanner

[Service]
ExecStart=clair -mode combo -conf /etc/clair/config.yaml

[Install]
WantedBy=multi-user.target
