# packages/c/clair/clair-watcher.service
[Unit]
Description=Clair Vulnerability Scanner

[Service]
ExecStart=clair -mode watcher -conf /etc/clair/config.yaml

[Install]
WantedBy=multi-user.target
