# packages/c/clair/clair-indexer.service
[Unit]
Description=Clair Vulnerability Scanner

[Service]
ExecStart=clair -mode indexer -conf /etc/clair/config.yaml

[Install]
WantedBy=multi-user.target
