# packages/p/ppp/modem@.service
[Unit]
Description=Modem /dev/ttyACM%i
BindsTo=dev-ttyACM%i.device
After=dev-ttyACM%i.device
Before=network.target
Wants=network.target

[Service]
Group=dialout
ExecStart=/usr/sbin/pppd call modem /dev/ttyACM%i

[Install]
WantedBy=multi-user.target
