# ndpmon/ndpmon.service
[Unit]
Description=IPv6 Neighbor Discovery Protocol Monitor
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/ndpmon

[Install]
WantedBy=multi-user.target

