# packages/o/oidentd/oidentd.socket
[Unit]
Conflicts=oidentd.service
Description=RFC 1413 compliant ident socket
Documentation=man:oidentd(8) man:oidentd.conf(5) man:oidentd_masq.conf(5)

[Socket]
Accept=true
ListenStream=113
PrivateDevices=true

[Install]
WantedBy=sockets.target