# packages/g/gnump3d/gnump3d.service
[Unit]
Description=GNU MP3 Streaming Server
Wants=network.target remote-fs.target
After=network.target remote-fs.target

[Service]
ExecStart=/usr/bin/gnump3d
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
