# packages/p/pptpd/pptpd.service
[Unit]
Description=PPTP Daemon
After=network.target

[Service]
ExecStart=/usr/sbin/pptpd --fg
StandardOutput=null

[Install]
WantedBy=multi-user.target