# packages/m/mycroft-core/mycroft-bus.service
[Unit]
Description=Mycroft bus service
After=network.target local-fs.target remote-fs.target
PartOf=mycroft.target

[Service]
Type=simple
PermissionsStartOnly=true
ExecStartPre=/usr/share/mycroft-core/create-initial-user-configuration.sh
ExecStart=/usr/bin/mycroft-messagebus
TimeoutStopSec=2s

[Install]
WantedBy=multi-user.target
