# packages/a/autossh/autossh.service
# DO NOT EDIT THIS FILE,
# use /etc/systemd/system/autossh@*.service.d/my.conf to override this file instead

[Unit]
Description=AutoSSH service for port 2222
After=network.target

[Service]
Environment="AUTOSSH_GATETIME=0"
ExecStart=/bin/echo "This line MUST be overridden.  See /usr/share/doc/packages/autossh/README.SUSE.md for details."

[Install]
WantedBy=multi-user.target

# DO NOT EDIT THIS FILE,
