# packages/d/dhcp/dhcpd.service
[Unit]
Description=ISC DHCPv4 Server
Before=multi-user.target
After=remote-fs.target network.target nss-lookup.target time-sync.target slapd.service sssd.service  ndsd.service

[Service]
Type=forking
Restart=on-abort
ExecStart=@LIBEXECDIR@/dhcp/dhcpd -4 start
ExecStop=@LIBEXECDIR@/dhcp/dhcpd -4 stop

[Install]
WantedBy=multi-user.target
Alias=dhcp-server.service

