# packages/d/dhcp/dhcrelay.service
[Unit]
Description=ISC DHCPv4 Relay Agent
Before=multi-user.target
After=remote-fs.target network.target nss-lookup.target time-sync.target ldap.service ndsd.service

[Service]
Type=forking
Restart=on-abort
ExecStart=@LIBEXECDIR@/dhcp/dhcrelay -4 start
ExecStop=@LIBEXECDIR@/dhcp/dhcrelay -4 stop
PIDFile=/run/dhcrelay.pid

[Install]
WantedBy=multi-user.target
Alias=dhcp-relay.service

