# packages/r/rpm/rpmconfigcheck.service
[Unit]
Description=Scan for unresolved .rpmnew, .rpmorig, and .rpmsave files
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/rpmconfigcheck

[Install]
WantedBy=default.target
