# packages/t/telnet/telnet.target
[Unit]
Description=Telnet target allowing to control multi setup
