# packages/i/iwlfwdump/50-iwlfwdump.rules
SUBSYSTEM=="devcoredump", ACTION=="add", RUN+="/sbin/iwlfwdump"
