# packages/a/apache2/apache2.target
[Unit]
Description=Apache target allowing to control multi setup

