# packages/g/goaccess/goaccess@.service
[Unit]
Description=goaccess instance: %i

[Service]
Type=oneshot
ExecStart=/usr/bin/goaccess --config-file=/etc/goaccess/%i.conf
