# packages/k/kbdsniffd/kbdsniffd.service
[Unit]
Description=kbdsniffd
After=remote-fs.target

[Service]
ExecStartPre=-/sbin/modprobe uinput
ExecStart=/usr/sbin/kbdsniffd

[Install]
WantedBy=multi-user.target

