# packages/o/openvpn/openvpn.target
[Unit]
Description=OpenVPN target allowing to start/stop all openvpn@.service instances at once

