# packages/o/openvpn/openvpn.service
[Unit]
Description=OpenVPN tunneling daemon instance using /etc/openvpn/%i.conf
After=network.target
PartOf=openvpn.target

[Service]
Type=notify
PrivateTmp=true
PIDFile=/run/openvpn/%i.pid
ExecStart=/usr/sbin/openvpn --daemon openvpn@%i --writepid /run/openvpn/%i.pid --cd /etc/openvpn/ --config %i.conf
# boo#1142830: "reload" does not work if openvpn drops root privileges after startup.
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target openvpn.target

