# packages/a/at/atd.service
[Unit]
Description=Deferred execution scheduler
After=nss-user-lookup.target time-sync.target

[Service]
ExecStart=/usr/sbin/atd -f

[Install]
WantedBy=multi-user.target
