# packages/g/golang-github-wrouesnel-postgres_exporter/prometheus-postgres_exporter.service
[Unit]
Description=Prometheus PostgreSQL exporter
After=postgresql.service
Wants=postgresql.service

[Service]
Restart=always
User=prometheus
EnvironmentFile=-/etc/sysconfig/prometheus-postgres_exporter
ExecStart=/usr/bin/prometheus-postgres_exporter $POSTGRES_EXPORTER_PARAMS
ExecReload=/bin/kill -HUP $MAINPID
TimeoutStopSec=20s

[Install]
WantedBy=multi-user.target
