# packages/s/sphinx/searchd.service
[Unit]
Description=Sphinx - SQL Full Text Search Engine
After=local-fs.target network.target mysqld.service postgresql.service

[Service]
User=sphinx
Group=sphinx
PIDFile=/run/sphinx/searchd.pid
ExecStart=/usr/bin/searchd --nodetach --config /etc/sphinx/sphinx.conf

[Install]
WantedBy=multi-user.target

