# packages/s/smartmontools/smartd_generate_opts.path
[Unit]
Description=Watch for changes in smartmontools sysconfig file
After=local-fs.target

[Path]
Unit=smartd_generate_opts.service
PathChanged=/etc/sysconfig/smartmontools

[Install]
WantedBy=default.target

