# packages/s/smartmontools/smartd_generate_opts.service
[Unit]
Description=Update smartd options
Wants=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/lib/smartmontools/generate_smartd_opts

[Install]
WantedBy=default.target

