# packages/p/pdnsd/pdnsd.service.in
[Unit]
Description=Caching DNS proxy (pdns)
Documentation=man:pdnsd(8)
After=network.target

[Service]
# BEGIN_NOT_IN_SLES12
# added automatically, for details please see
# https://en.opensuse.org/openSUSE:Security_Features#Systemd_hardening_effort
ProtectSystem=full
ProtectHome=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
# end of automatic additions 
# END_NOT_IN_SLES12
Type=simple
EnvironmentFile=/etc/sysconfig/pdnsd
ExecStart=/usr/sbin/pdnsd $PDNSD_ARGS

[Install]
WantedBy=multi-user.target
