# packages/p/pmacct/pmacct.sfacctd.service
[Unit]
Description=sflow accounting daemon
After=network.target

[Service]
Type=forking
EnvironmentFile=-/etc/sysconfig/sfacctd
ExecStart=/usr/sbin/sfacctd -f ${SFACCTD_CONF} $DAEMON_OPTS

[Install]
WantedBy=multi-user.target
