# packages/r/redis/redis-sentinel.target
[Unit]
Description=Redis Sentinel target allowing to start/stop all redis-sentinel@.service instances at once

