# packages/f/fwknop/fwknopd.service
[Unit]
Description=Firewall Knock Operator Daemon
After=network.target

[Service]
Type=forking
PIDFile=/var/run/fwknopd.pid
ExecStart=/usr/sbin/fwknopd
ExecReload=/usr/sbin/fwknopd -R
ExecStop=/usr/sbin/fwknopd -K

[Install]
WantedBy=multi-user.target
