# packages/i/icecream/icecc-scheduler.service.in
[Unit]
Description=Icecream distributed compiler scheduler

[Service]
# added automatically, for details please see
# https://en.opensuse.org/openSUSE:Security_Features#Systemd_hardening_effort
ProtectSystem=full
ProtectHome=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
# end of automatic additions 
Type=simple
EnvironmentFile=-/etc/sysconfig/icecream
User=icecream
Group=icecream
SyslogIdentifier=icecc-scheduler
ExecStart=@LIBEXECDIR@/icecc/icecc-scheduler-wrapper -u icecream

[Install]
WantedBy=multi-user.target
