# packages/i/icecream/iceccd.service.in
[Unit]
Description=Icecream Distributed Compiler
After=network.target nss-lookup.target

[Service]
# added automatically, for details please see
# https://en.opensuse.org/openSUSE:Security_Features#Systemd_hardening_effort
ProtectSystem=full
ProtectHome=true
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
# end of automatic additions 
Type=simple
EnvironmentFile=-/etc/sysconfig/icecream
SyslogIdentifier=iceccd
ExecStart=@LIBEXECDIR@/icecc/iceccd-wrapper -u icecream -b /var/cache/icecream/
Nice=5

[Install]
WantedBy=multi-user.target
