# packages/i/icecream/icecc-scheduler.service
[Unit]
Description=Icecream distributed compiler scheduler

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/icecream
User=icecream
Group=icecream
SyslogIdentifier=icecc-scheduler
ExecStart=/usr/lib/icecc/icecc-scheduler-wrapper -u icecream

[Install]
WantedBy=multi-user.target
