# packages/x/xdm/xdm.service
[Unit]
Description=X Display Manager
Documentation=man:xdm(8)
Conflicts=getty@tty7.service plymouth-quit.service
After=ypbind.service gpm.service winbind.service acpid.service
Wants=systemd-user-sessions.service systemd-logind.service
After=systemd-user-sessions.service systemd-logind.service
OnFailure=plymouth-quit.service

[Service]
ExecStart=/usr/bin/xdm -nodaemon

[Install]
Alias=display-manager.service
