# packages/q/qos/qos.service.in
[Unit]
Description=QOS traffic shaping to avoid bufferbloat
After=network-pre.target
Documentation=https://github.com/archiecobbs/qos

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@qosprog@ start
ExecStop=@qosprog@ stop
