# packages/n/neard/neard.service
[Unit]
Description=neard service
Documentation=man:neard(8)

[Service]
Type=dbus
BusName=org.neard
ExecStart=/usr/lib/nfc/neard -n
NotifyAccess=main
#WatchdogSec=10
#Restart=on-failure
#CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_BIND_SERVICE CAP_NET_RAW
LimitNPROC=1

[Install]
Alias=dbus-org.neard.service
