# packages/i/incus/incus-user.socket
[Unit]
Description=Incus Container Hypervisor (User Daemon) Socket

[Socket]
ListenStream=/run/incus/unix.socket.user
SocketGroup=incus
SocketMode=0660
Service=incus-user.service

[Install]
WantedBy=sockets.target
