# canna/canna.service
[Unit]
Description=Canna Japanese Conversion Engine

[Service]
Type=forking
ExecStartPre=/bin/rm -rf /tmp/.iroha_unix
ExecStart=/usr/sbin/cannaserver -s -u wnn -r /var/lib/canna
ExecStopPost=/bin/rm -rf /tmp/.iroha_unix

[Install]
WantedBy=multi-user.target
