# packages/r/rkhunter/rkhunter.service
[Unit]
Description=rkhunter sevice
Wants=rkhunter.timer

[Service]
Type=oneshot
ExecStart=/usr/bin/rkhunter --propupd
