# packages/r/rkhunter/rkhunter.timer
[Unit]
Description=rkhunter timer

[Timer]
Unit=rkhunter.service
OnBootSec=15m
OnUnitActiveSec=2d

[Install]
WantedBy=timers.target
