# packages/o/openwsman/openwsman.service
[Unit]
Description=Openwsman WS-Management Service
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/sbin/openwsmand -S

[Install]
WantedBy=multi-user.target
