# packages/s/squid/squid.service
[Unit]
Description=Squid caching proxy
Documentation=man:squid(8)
After=network.target named.service nss-lookup.service

[Service]
Type=forking
ExecStartPre=%{_libdir}/squid/initialize_cache_if_needed.sh
ExecStart=/usr/sbin/squid -FC
ExecReload=/usr/bin/kill -HUP $MAINPID
LimitNOFILE=4096
PIDFile=/run/squid.pid

[Install]
WantedBy=multi-user.target

