# packages/s/SUSEConnect/suseconnect-keepalive.timer
[Unit]
Description=Schedule daily SUSEConnect --keepalive calls
After=network-online.target

[Timer]
# Run this timer every day at a randomized 24h delay.
OnCalendar=daily
RandomizedDelaySec=24h

[Install]
WantedBy=timers.target
