# pure-ftpd/pure-ftpd.service
[Unit]
Description=Pure-FTPd FTP server
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/pure-ftpd /etc/pure-ftpd/pure-ftpd.conf

[Install]
WantedBy=multi-user.target

