# fwnn/fwnn.service
[Unit]
ConditionPathExists=/etc/FreeWnn/ja/jserverrc
Description=Free Wnn Kanji Server used for input of Japanese

[Service]
Type=forking
User=wnn
ExecStartPre=/bin/rm -f /tmp/jd_sockV4
ExecStart=/usr/bin/jserver
ExecStopPost=/bin/rm -f /tmp/jd_sockV4

[Install]
WantedBy=multi-user.target
