# packages/s/spamassassin/sa-update.service
[Unit]
Description=Update spamassassin
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/usr/sbin/timed-sa-update
