# packages/d/data-partition-service/data_part.service
[Unit]
Description=Add a data partition
Requires=systemd-fsck-root.service
After=systemd-fsck-root.service
Before=systemd-remount-fs.service
Wants=local-fs-pre.target
ConditionPathExists=/etc/data_part.cfg

# Required to let the service run in early stages (before local-fs.target)
DefaultDependencies=no
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/share/data-partition-service/create_data_part

[Install]
WantedBy=local-fs.target
