# packages/i/ipsec-tools/racoon.service
[Unit]
Description=Racoon IKEv1 key management daemon for IPSEC
After=network.target racoon-setkey.service
ConditionPathExists=/etc/racoon/racoon.conf

[Service]
EnvironmentFile=-/etc/sysconfig/racoon
ExecStart=/usr/sbin/racoon -F $RACOON_OPTIONS -f /etc/racoon/racoon.conf
ExecReload=/bin/kill -HUP $MAINPID

[Install]
Also=racoon-setkey.service
WantedBy=multi-user.target
