# packages/l/loki/loki.target
[Unit]
Description=Target to restart all parts of Grafana Loki

[Install]
WantedBy=multi-user.target
