# packages/d/distcc/distccd.service
[Unit]
Description=Distccd: A Distributed Compilation Server
After=network.target

[Service]
EnvironmentFile=-/etc/sysconfig/distccd
ExecStart=/usr/bin/distccd --verbose --no-detach --daemon --port 3632 -N 15 --allow $ALLOWED_SERVERS $EXTRA_PARAMS

[Install]
WantedBy=multi-user.target
