# packages/s/spacenavd/spacenavd.service
[Unit]
Description=Userspace Daemon of the spacenav driver.

[Service]
Type=forking
PIDFile=/run/spnavd.pid
ExecStart=/usr/sbin/spacenavd

[Install]
WantedBy=multi-user.target
