# packages/o/openSUSE-build-key/openSUSE-build-key-import.timer
[Unit]
Description=Timer that starts importing new GPG keys into RPM Database.

[Timer]
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
