# packages/k/kdebase4-workspace/ksysguardd.service
[Unit]
Description=KDE ksysguard daemon

[Service]
Type=forking
PIDFile=/var/run/ksysguardd.pid
ExecStart=/usr/bin/ksysguardd -d
PrivateTmp=true

[Install]
WantedBy=multi-user.target