# packages/o/openslp/slpd.service
[Unit]
Description=OpenSLP daemon for the Service Location Protocol
After=network.target nss-lookup.target

[Service]
Type=notify
ExecStart=/usr/sbin/slpd -d
ExecReload=/bin/kill -HUP $MAINPID
[Install]
WantedBy=multi-user.target