# packages/l/linuxptp/phc2sys.service
[Unit]
Description=PTP: Synchronize two clocks
After=ntpdate.service

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/phc2sys
ExecStart=/usr/sbin/phc2sys $OPTIONS

[Install]
WantedBy=multi-user.target
