# ypserv/yppasswdd.service
[Unit]
Description=NIS/YP Users Passwords Change Server
Requires=rpcbind.service ypserv.service
After=network.target rpcbind.service ypserv.service

[Service]
Type=notify
NotifyAccess=all
EnvironmentFile=-/etc/sysconfig/ypserv
ExecStart=/usr/lib/yp/yppasswdd-systemd-exec

[Install]
WantedBy=multi-user.target
