# packages/v/vpnc/vpnc.service
[Unit]
Description=VPNC connection to %i
Documentation=man:vpnc(8) http://www.unix-ag.uni-kl.de/~massar/vpnc/
Requires=network.target
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/vpnc --pid-file=/run/vpnc@%i.pid /etc/vpnc/%i.conf
PIDFile=/run/vpnc@%i.pid

[Install]
WantedBy=multi-user.target
