# packages/c/courier-imap/courier-imap-gencert.service
[Unit]
Description=Courier IMAP (SSL) Cert Generation
ConditionPathExists=!/etc/ssl/private/imapd.pem

[Service]
Type=oneshot
ExecStart=/usr/sbin/mkimapdcert
RemainAfterExit=yes 
