# packages/e/ebtables/ebtables.service
[Unit]
Description=Ethernet Bridge Filtering tables

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/ebtables start
ExecStop=/usr/lib/ebtables stop

[Install]
WantedBy=multi-user.target
