# packages/r/rp-pppoe/pppoe-server.service
[Unit]
Description=PPPoE Server.
After=syslog.target

[Service]
Group=dialout
ExecStart=/usr/sbin/pppoe-server

[Install]
WantedBy=multi-user.target
