# packages/r/rp-pppoe/pppoe.service
[Unit]
Description=PPPoE
After=network.target

[Service]
Group=dialout
Type=oneshot
RemainAfterExit=true
ExecStart=-/usr/sbin/pppoe-start
ExecStop=/usr/sbin/pppoe-stop

[Install]
WantedBy=multi-user.target
