# packages/i/inn/inn.service
[Unit]
Description=Inter Net News Server (INN)
After=network-online.target
After=nss-lookup.target

[Service]
User=news
PIDFile=/run/news/innd.pid
Type=forking
ExecStart=@LIBEXECDIR@/news/bin/rc.news
ExecStop=@LIBEXECDIR@/news/bin/rc.news stop
ExecReload=@LIBEXECDIR@/news/bin/rc.news reload
KillMode=control-group

[Install]
WantedBy=multi-user.target
