# packages/f/fdo-client/fdoclient.service
[Unit]
Description=FDO client
After=remote-fs.target network-online.target 
Wants=network-online.target

[Service]
Type=simple
Environment=TERM=linux
ExecStart=/usr/bin/fdo-client-service
RemainAfterExit=no
TimeoutSec=0

[Install]
WantedBy=default.target

