# govpn/govpn.target
[Unit]
Description=GoVPN target allowing to start/stop all govpn@.service instances at once
