# packages/s/s390-tools/appldata.service
[Unit]
Description=Linux - z/VM Monitor Stream
After=network-online.target remote-fs.target 
Wants=network-online.target remote-fs.target
ConditionPathExists=/proc/sys/appldata/interval
ConditionPathExists=!/var/lock/appldata

[Service]
Type=oneshot
RemainAfterExit=yes

ExecStart=/usr/lib/systemd/scripts/appldata start
ExecStartPost=/usr/lib/systemd/scripts/appldata status
ExecStop=/usr/lib/systemd/scripts/appldata stop

[Install]
WantedBy=default.target
