# packages/s/s390-tools/hsnc.service
[Unit]
Description=Start the qeth HiperSockets Network Concentrator
After=network-online.target remote-fs.target 
Wants=network-online.target remote-fs.target
ConditionPathExists=/sys/devices/qeth

[Service]
Type=oneshot
RemainAfterExit=yes

ExecStart=/usr/lib/systemd/scripts/hsnc start
ExecStartPost=/usr/lib/systemd/scripts/hsnc status
ExecStop=/usr/lib/systemd/scripts/hsnc stop

[Install]
WantedBy=default.target
