# packages/s/s390-tools/virtsetup.service
[Unit]
Description=Perform hypervisor-specific setup and cleanup tasks
DefaultDependencies=no
Wants=default.target
After=default.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/scripts/virtsetup.sh

[Install]
WantedBy=default.target
