# packages/h/hostname/nis-domainname.service.in
[Unit]
Description=Read and set NIS domainname from /etc/sysconfig/network/config
Before=ypbind.service yppasswdd.service ypserv.service ypxfrd.service sysinit.target
DefaultDependencies=no
Conflicts=shutdown.target
ConditionPathExists=/etc/sysconfig/network/config

[Service]
ExecStart=@LIBEXECDIR@/hostname/nis-domainname
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
