# packages/s/shadowsocks-libev/shadowsocks-libev-redir@.service
[Unit]
Description=Daemon to start Shadowsocks-libev-redir
Wants=network-online.target
After=network.target

[Service]
Type=forking
PIDFile=/var/run/shadowsocks-libev-redir@%i.pid
ExecStart=/usr/bin/ss-redir -c /etc/shadowsocks/%i.json -f /var/run/shadowsocks-libev-redir@%i.pid -u --fast-open
Restart=on-failure

[Install]
WantedBy=multi-user.target
