# packages/p/pgbackrest/pgbackrest.target
[Unit]
Description=Target to restart/stop all parts of pgbackrest

