# packages/p/pgbackrest/pgbackrest-diff@.service
[Unit]
Description=pgbackrest backup diff for %i
Documentation=https://pgbackrest.org/configuration.html
Documentation=man:pgbackrest(1) man:pgbackrest.conf(5)
ConditionACPower=true
PartOf=pgbackrest.target

[Service]
Type=oneshot
User=postgres
Group=postgres
ExecStart=/usr/bin/pgbackrest --log-level-console=info --stanza=%i --type=diff backup
# Give a reasonable amount of time for job to finish
TimeoutSec=7200

