# packages/p/pgbackrest/pgbackrest-incr@.timer
[Unit]
Description=Daily pgbackrest incremental backup (mon-sat)
Documentation=https://pgbackrest.org/configuration.html
Documentation=man:pgbackrest(1) man:pgbackrest.conf(5)
PartOf=pgbackrest.target

[Timer]
OnCalendar=Mon..Sat *-*-* 05:05:05
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
