# packages/p/pgbackrest/pgbackrest-full@.timer
[Unit]
Description=Monthly pgbackrest full backup (1st sun)
Documentation=https://pgbackrest.org/configuration.html
Documentation=man:pgbackrest(1) man:pgbackrest.conf(5)
PartOf=pgbackrest.target

[Timer]
#OnCalendar=Sun *-*-1..7 05:05:05
OnCalendar=Sun *-*-* 05:05:05
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
