# packages/p/pgbackrest/pgbackrest.service
[Unit]
Description=pgBackRest Server
Documentation=https://pgbackrest.org/configuration.html
Documentation=man:pgbackrest(1) man:pgbackrest.conf(5)
After=network.target
StartLimitIntervalSec=0
PartOf=pgbackrest.target

[Service]
Type=simple
Restart=always
RestartSec=1
User=postgres
ExecStart=/usr/bin/pgbackrest server
ExecReload=kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
