# packages/i/ipmiutil/ipmiutil_checksel.timer
[Unit]
Description=Automatic execution of ipmiutil's checksel

[Timer]
OnCalendar=Daily
Persistent=true

[Install]
WantedBy=timers.target