# packages/y/ypserv/yppasswdd.service
[Unit]
Description=NIS/YP Users Passwords Change Server
Requires=rpcbind.service ypserv.service
After=network.target rpcbind.service ypserv.service

[Service]
Type=notify
NotifyAccess=all
EnvironmentFile=-/usr/etc/default/yppasswdd
EnvironmentFile=-/etc/default/yppasswdd
ExecStart=/usr/libexec/yp/yppasswdd-systemd-exec

[Install]
WantedBy=multi-user.target
