# packages/y/ypserv/ypserv.service
[Unit]
Description=NIS/YP Server
Requires=rpcbind.service
After=network.target rpcbind.service slpd.service

[Service]
Type=notify
NotifyAccess=all
EnvironmentFile=-/etc/sysconfig/ypserv
ExecStart=/usr/sbin/ypserv -f $YPSERV_ARGS

[Install]
WantedBy=multi-user.target
