# packages/y/ypserv/ypxfrd.service
[Unit]
Description=NIS/YP Maps Transferring Accelerator
Requires=rpcbind.service ypserv.service
After=network.target rpcbind.service ypserv.service

[Service]
Type=notify
NotifyAccess=all
EnvironmentFile=-/etc/sysconfig/ypserv
ExecStart=/usr/sbin/rpc.ypxfrd -f $YPXFRD_ARGS

[Install]
WantedBy=multi-user.target
