# packages/m/mariadb/mariadb.target
[Unit]
Description=MySQL target allowing to control multi setup
