# packages/d/dehydrated/dehydrated.timer
[Unit]
Description=Run Certificate Update Runner for Dehydrated 

[Timer]
OnCalendar=daily
# Two hour window
RandomizedDelaySec=7200
 
[Install]
WantedBy=timers.target
