# packages/h/hp-drive-guard/hp-drive-guard.service
[Unit]
Description=HP DriveGuard Daemon

[Service]
Type=dbus
BusName=com.hp.DriveGuard
ExecStart=/usr/sbin/hp-drive-guard-daemon
