# packages/s/ser2net/ser2net.service
[Unit]
Description=Proxy that allows tcp connections to serial ports
After=network.target

[Service]
ExecStart=/usr/sbin/ser2net -n

[Install]
WantedBy=multi-user.target
