# rng-tools/rng-tools.service
[Unit]
Description=Start the rngd daemon
ConditionVirtualization=!container
Conflicts=haveged.service

[Service]
CapabilityBoundingSet=CAP_SYS_ADMIN
ExecStart=/usr/sbin/rngd --foreground --fill-watermark=3700

[Install]
WantedBy=multi-user.target
Alias=rngd.service
