# packages/g/geoipupdate/geoipupdate.timer
[Unit]
Description=Scheduled update of GeoIP2/GeoLite2 databases
Documentation=man:geoipupdate(1) man:GeoIP.conf(5)

[Timer]
OnCalendar=Wed,Sat America/New_York

AccuracySec=1h
RandomizedDelaySec=3h

Persistent=true

[Install]
WantedBy=timers.target
