# packages/g/geoipupdate/geoipupdate.service
[Unit]
Description=Update GeoIP databases
Documentation=man:geoipupdate(1) man:GeoIP.conf(5)

[Service]
Type=oneshot
WorkingDirectory=/var/lib/GeoIP

ExecStart=/usr/bin/geoipupdate
