# packages/g/geoipupdate/geoipupdate.timer
[Unit]
Description=Weekly update of GeoLite2 databases
Documentation=man:geoipupdate(1) man:GeoIP.conf(5)

[Timer]
OnUnitInactiveSec=1w
OnStartupSec=1h
AccuracySec=1h

[Install]
WantedBy=timers.target

