# packages/l/logdigest/logdigest.timer
[Unit]
Description=Daily execution of logdigest

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
