# packages/m/mednafen-server/mednafen.service
[Unit]
Description=mednafen server
After=network.target

[Service]
User=mednafen
Group=games
ExecStart=/usr/bin/mednafen-server /var/lib/mednafen/mednafen-server.conf

[Install]
WantedBy=multi-user.target
