# google-compute-engine/google-optimize-local-ssd.service
[Unit]
Description=Google Compute Engine Local SSD Optimization
After=local-fs.target
Wants=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/google_optimize_local_ssd
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
