# packages/x/xfs/xfs.service
[Unit]
Documentation=man:xfs(1)
SourcePath=/usr/sbin/rcxfs
Description=LSB: X Font Server
After=network-online.target
After=nss-lookup.target
After=remote-fs.target
Wants=network-online.target
Wants=nss-lookup.target
Wants=remote-fs.target

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
SuccessExitStatus=5 6
ExecStart=/usr/sbin/rcxfs start
ExecStop=/usr/sbin/rcxfs stop
ExecReload=/usr/sbin/rcxfs reload
