# packages/m/mailgraph/mailgraph.service
[Unit]
Description=RRDtool frontend for Mail statistic (mailgraph)
Documentation=http://mailgraph.schweikert.ch/
After=mail-transfer-agent.target

[Service]
EnvironmentFile=-/etc/sysconfig/mailgraph
Nice=19
ExecStart=/usr/bin/mailgraph.pl $MAILGRAPH_OPTS --logfile ${MAILGRAPH_LOG_FILE} --logtype ${MAILGRAPH_LOG_TYPE} --daemon-rrd ${MAILGRAPH_RRD} --daemon-log ${MAILGRAPH_LOG}

[Install]
WantedBy=multi-user.target
