# packages/p/ppc64-diag/ppc64-diag-encl.timer
[Unit]
Description=Daily execution of ppc64-diag's enclosure diagnostic

[Timer]
OnCalendar=Daily
Persistent=true

[Install]
WantedBy=timers.target
