# packages/k/kafka/kafka-zookeeper.service
[Unit]
Description=Single node Zookeeper cluster for local Kafka testing. Not suitable for production use.
After=syslog.target network.target

[Service]
Type=simple
User=kafka
Group=kafka
WorkingDirectory=/run/kafka-zookeeper
Environment='LOG_DIR=/var/log/kafka'
PrivateTmp=yes
ExecStart=/usr/bin/kafka-zookeeper-server-start.sh /etc/kafka/zookeeper.properties

[Install]
WantedBy=multi-user.target
