# qemu/qemu-ga@.service
[Unit]
Description=QEMU Guest Agent
Documentation=http://wiki.qemu.org/Features/GuestAgent
BindsTo=dev-%i.device
After=dev-%i.device

[Service]
Type=simple
ExecStart=-/usr/bin/qemu-ga -p /dev/%I
Restart=always
RestartSec=0

[Install]
WantedBy=dev-%i.device
