# packages/t/tigervnc/xvnc.socket
[Unit]
Description=Xvnc Server
Before=xvnc.target
Wants=xvnc.target

[Socket]
ListenStream=5901
Accept=yes

[Install]
WantedBy=sockets.target
