# samplicator/samplicator.service
[Unit]
Description=Send copies of (UDP) datagrams to multiple receivers, with optional sampling and spoofing
After=network.target

[Service]
Type=forking
EnvironmentFile=/etc/sysconfig/samplicator
ExecStart=/usr/sbin/samplicate -c /etc/samplicator.conf -d 0 -f $SAMPLICATOR_OPTIONS=""
Restart=on-abort
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
