# packages/i/i4l-base/45-isdn.rules
# do not edit this file, it will be overwritten on update

# CAPI devices
SUBSYSTEM=="capi", KERNEL=="capi", NAME="capi20", GROUP="dialout"
SUBSYSTEM=="tty", KERNEL=="capi[0-9]*", NAME="capi/%n"

# ISDN devices
KERNEL=="isdn*", MODE="600"
KERNEL=="isdnctrl*", MODE="600"
KERNEL=="ippp*", MODE="600"
KERNEL=="isdninfo*", MODE="444"

# ISDN USB
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/c00*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/1000*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/1900*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/2000*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/2200*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/2300*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/2800*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/3000*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/3500*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/3600*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/2805*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/2806*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/4401*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="57c/4601*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="483/4810*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="959/2bd0*", RUN+="isdn.sh"
SUBSYSTEM=="usb", ENV{PRODUCT}=="9bf/*", RUN+="isdn.sh"

# ISDN PCMCIA
# AVM B1 models
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa95D42008pb845DC335pc*pd*", RUN+="isdn.sh b1 1"
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa95D42008pb81E10430pc*pd*", RUN+="isdn.sh b1 2"
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa95D42008pb18E8558Apc*pd*", RUN+="isdn.sh b1 3"
# AVM A1 models (aka Fritz!Card PCMCIA)
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa95D42008pbADC9D4BBpc*pd*", RUN+="isdn.sh a1 1"
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa8D9761C8pb01C5AA7Bpc*pd*", RUN+="isdn.sh a1 2"
# teles
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa67B50EAEpbE9E70119pc*pd*", RUN+="isdn.sh teles 1"
# sedlbauer
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa81FB79F5pbF3612E1Dpc6B95C78Apd*", RUN+="isdn.sh sedlbauer 1"
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa81FB79F5pbE4E9BC12pc397B7E90pd*", RUN+="isdn.sh sedlbauer 2"
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa81FB79F5pbE4E9BC12pc2E5C7FCEpd*", RUN+="isdn.sh sedlbauer 3"
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa81FB79F5pbE4E9BC12pc8DB143FEpd*", RUN+="isdn.sh sedlbauer 4"
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa81FB79F5pbE4E9BC12pcB391AB4Cpd*", RUN+="isdn.sh sedlbauer 5"
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*paD79E0B84pb21D083AEpc*pd*", RUN+="isdn.sh sedlbauer 6"
# elsa
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa983DE2C4pb333BA257pc*pd*", RUN+="isdn.sh elsa 1"
SUBSYSTEM=="pcmcia", ENV{MODALIAS}=="pcmcia:m*c*f*fn*pfn*pa639E5718pb333BA257pc*pd*", RUN+="isdn.sh elsa 2"
