# packages/m/mt-st/61-storage-tape-init.rules
# Initialize SCSI magnetic tape drives

ACTION!="remove", SUBSYSTEM=="scsi_tape", KERNEL=="st*[0-9]", RUN+="/usr/sbin/stinit %n"
